# ahoydtu-postgres-fetcher.cron
*/1 * * * * root /usr/local/bin/ahoydtu-postgres-fetcher.py
